# templates/etc/init.d/crowd.erb
#!/bin/sh
#
# Crowd startup script
#
# chkconfig:   2345 80 05
#
# description: Atlassian Crowd

### BEGIN INIT INFO
# Provides:
# Required-Start:
# Required-Stop:
# Should-Start:
# Should-Stop:
# Default-Start:
# Default-Stop:
# Short-Description:
# Description:
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions


RUN_AS_USER=<%= scope.lookupvar('crowd::user') %>
CATALINA_HOME=<%= scope.lookupvar('crowd::webappdir') %>/apache-tomcat
exec="${CATALINA_HOME}/bin/startup.sh"
pidfile="${CATALINA_HOME}/work/crowd.pid"

start() {
  [ -x $exec ] || exit 5
  echo -n $"Starting Crowd: "
  # if not running, start it up here, usually something like "daemon $exec"
  daemon --user "${RUN_AS_USER}" "${exec}"
  retval=$?
  echo
  return $retval
}

stop() {
  echo -n $"Stopping Crowd: "
  # stop it here, often "killproc $prog"
  daemon --user "${RUN_AS_USER}" "$CATALINA_HOME/bin/shutdown.sh"
  retval=$?
  echo
  return $retval
}
status() {
  if [ -f $pidfile ]; then
    PID=$(cat $pidfile)
    if [ -z "`ps axf | grep ${PID} | grep -v grep`" ]; then
      echo "Crowd: Process is dead but PID file exists at $pidfile"
      exit 2
    else
      echo "Crowd is running with pid ${PID}..."
      exit 0
    fi
  else
    echo "Crowd is not running"
    exit 1
  fi
}

case "$1" in
  start)
    start
    ;;
  stop)
    stop
    ;;
  status)
    status
    ;;
  restart)
    stop
    sleep 10
    start
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart}"
esac

exit 0
