# templates/crowd.init.erb
#!/bin/sh
#
# Crowd startup script
#
# chkconfig:   2345 80 05
#
# description: Atlassian Crowd

### BEGIN INIT INFO
# Provides: crowd
# Required-Start: $remote_fs $syslog
# Required-Stop: $remote_fs $syslog
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: Start and stop the Atlassian Crowd service
# Description: Provides the crowd authentication service
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

RUN_AS_USER="<%= scope.lookupvar('crowd::user') %>"
CATALINA_HOME="<%= scope.lookupvar('crowd::app_dir') %>/apache-tomcat"

crowd_start() {
  echo "Starting Crowd: "
  if [ "x$USER" != "x$RUN_AS_USER" ]; then
    sudo -u $RUN_AS_USER "$CATALINA_HOME/bin/startup.sh"
  else
    $CATALINA_HOME/bin/startup.sh
  fi
  echo "done."
}

crowd_stop() {
  echo "Shutting down Crowd: "
  if [ "x$USER" != "x$RUN_AS_USER" ]; then
    sudo -u $RUN_AS_USER "$CATALINA_HOME/bin/shutdown.sh"
  else
    $CATALINA_HOME/bin/shutdown.sh
  fi
  echo "done."
}

crowd_status() {
  status -p $pidfile crowd
  retval=$?
  echo
  return $retval
}

case "$1" in
  start)
    crowd_start
    ;;
  stop)
    crowd_stop
    ;;
  status)
    crowd_status
    ;;
  restart)
    crowd_stop
    sleep 10
    crowd_start
    ;;
  *)
    echo "Usage: $0 {start|stop|status|restart}"
esac
